# roles/base/files/hobbit-backup.timer
[Unit]
Description=Run Hobbit Backup Weekly

[Timer]
OnCalendar=Sun 02:00
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
