# files/hobbit-bridge.service.j2
[Unit]
Description=Hobbit Control Bridge
After=network.target

[Service]
Type=simple
User=hobbit
WorkingDirectory=/home/hobbit/hobbit/bridge
ExecStart=/usr/bin/node bridge.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=GAMING_PC_HOST=192.168.0.69
Environment=SPOTIFY_CLIENT_ID={{ spotify_client_id }}
Environment=SPOTIFY_CLIENT_SECRET={{ spotify_client_secret }}
Environment=SPOTIFY_REDIRECT_URI=https://{{ tailscale_fqdn }}/api/control/spotify/callback
Environment=CAMERA_PASSWORD={{ camera_password }}
Environment="WIFI_SSID={{ wifi_ssid }}"
Environment="WIFI_PASSWORD={{ wifi_password }}"

[Install]
WantedBy=multi-user.target
